# services/nvm-update.service
[Unit]
Description=Update Node.js (nvm) and global npm tools
Documentation=https://github.com/nvm-sh/nvm
Wants=network-online.target
After=network-online.target
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=%h/.local/bin/nvm-update.sh

Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=%h
Environment=NVM_DIR=%h/.nvm
Environment=NVM_NODE_ALIAS=default
Environment=NVM_NODE_MAJOR=22
Environment="NVM_GLOBAL_TOOLS=npm typescript yarn"
Environment="AI_TOOLS_GLOBAL_TOOLS=npm @anthropic-ai/claude-code"

# NoNewPrivileges is intentionally absent: nvm-update.sh calls
# "sudo -u ai-tools" to update /opt/ai-tools, which requires the setuid bit.
PrivateTmp=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nvm-update
