# src/usr/lib/systemd/system/ai-tools-handback.socket
# /usr/lib/systemd/system/ai-tools-handback.socket
# Listening socket for the ai-tools handback privilege bridge.
#
# DAC gate: SocketMode=0660 + SocketGroup=@SANDBOX_GROUP@ so only root and the
# sandbox group can connect.  (The sandbox user @SANDBOX_USER@ is a member of
# @SANDBOX_GROUP@; ordinary users are not.)  The handler performs a second
# AuthN gate via SO_PEERCRED.
#
# RuntimeDirectory creates /run/ai-tools at socket startup with mode 0750
# root:@SANDBOX_GROUP@, so group members can traverse the directory to reach
# the socket file name but cannot list other entries.
#
# Accept=yes: one ai-tools-handback@.service instance per connection; each
# serves exactly one request and exits.  The listening socket (and /run/ai-tools)
# are destroyed when this unit stops.

[Unit]
Description=ai-tools handback privilege bridge socket
Documentation=file:/usr/local/sbin/ai-tools/ai-tools-handback

[Socket]
ListenStream=/run/ai-tools/handback.sock
SocketMode=0660
SocketGroup=@SANDBOX_GROUP@
Accept=yes
# Cap simultaneous connections from SANDBOX_USER.  Each accepted connection
# spawns a handler process and a helper subprocess (both as root); without a
# limit a flood of connections is an unbounded-process-count DoS.  20 is well
# above the expected hook concurrency (one PostToolUse + one Stop per turn).
MaxConnections=20
# Trigger rate limit made explicit.  An Accept=yes socket that is triggered more
# than TriggerLimitBurst times within TriggerLimitIntervalSec is DISABLED until a
# manual `systemctl restart`.  Tripping it is fail-SAFE (the bridge stops, no
# privileged action runs, the agent's files simply stay SANDBOX_USER-owned and the
# hooks no-op via `|| true`) -- but it would still halt all handback, so the burst
# is set comfortably above any legitimate session's hook rate while still bounding
# a flood.  Raise TriggerLimitBurst if a heavy interactive session ever trips it.
TriggerLimitIntervalSec=10s
TriggerLimitBurst=200
RuntimeDirectory=ai-tools
RuntimeDirectoryMode=0750
RuntimeDirectoryGroup=@SANDBOX_GROUP@

[Install]
WantedBy=sockets.target
