# src/usr/lib/systemd/system/ai-tools-handback.socket
# /usr/lib/systemd/system/ai-tools-handback.socket
# Listening socket for the ai-tools handback privilege bridge.
#
# DAC gate: SocketMode=0660 + SocketGroup=@SANDBOX_GROUP@ so only root and the
# sandbox group can connect.  (The sandbox user @SANDBOX_USER@ is a member of
# @SANDBOX_GROUP@; ordinary users are not.)  The handler performs a second
# AuthN gate via SO_PEERCRED.
#
# RuntimeDirectory creates /run/ai-tools root:root on every socket (re)start.  Its
# GROUP cannot be set declaratively on the systemd 252 baseline (RHEL/Rocky 9):
# RuntimeDirectoryGroup= first exists in v256, and on 252 it is parsed-and-ignored
# ("Unknown key name 'RuntimeDirectoryGroup' ... ignoring").  A chgrp in ExecStartPost
# does not stick either -- RuntimeDirectory= re-asserts root:root on each start (the
# dir stayed root:root with zero AVCs even under `semodule -DB`).  So directory access
# is governed by MODE, not group: 0711 makes it world-traversable but not listable, so
# a @SANDBOX_GROUP@ member can resolve the socket name to connect() while the directory
# contents stay unlistable.  This is not the access boundary -- the socket file is
# 0660 root:@SANDBOX_GROUP@ (SocketMode/SocketGroup) and the daemon re-checks the peer
# uid via SO_PEERCRED, so a non-group uid that traverses in still gets EACCES.
#
# ExecStartPost restorecon: /run is tmpfs, so RuntimeDirectory tears the path
# down and recreates it on every socket (re)start -- boot, manual restart, or
# automatic recovery from a tripped TriggerLimit -- deriving its label from
# whatever context PID1 currently resolves for the path. That resolution has
# proven to drift to var_run_t (a stale cached file_contexts lookup; see
# install-selinux.sh _relabel_runtime), which then denies ai_tools_t write on
# the socket. Re-running install-selinux.sh's relabel step fixes the path that
# exists at that moment, but the fix does not survive the next recreation.
# restorecon here makes the correct ai_tools_run_t label a property of every
# start instead of a one-time external repair: it runs only when this unit
# (re)starts -- not per connection -- and is a cheap, idempotent no-op once the
# label already matches policy, so it adds no meaningful overhead.
#
# Accept=yes: one ai-tools-handback@.service instance per connection; each
# serves exactly one request and exits.  The listening socket (and /run/ai-tools)
# are destroyed when this unit stops.

[Unit]
Description=ai-tools handback privilege bridge socket
Documentation=file:/usr/local/sbin/ai-tools/ai-tools-handback

[Socket]
ListenStream=/run/ai-tools/handback.sock
SocketMode=0660
SocketGroup=@SANDBOX_GROUP@
Accept=yes
# Cap simultaneous connections from SANDBOX_USER.  Each accepted connection
# spawns a handler process and a helper subprocess (both as root); without a
# limit a flood of connections is an unbounded-process-count DoS.  20 is well
# above the expected hook concurrency (one PostToolUse + one Stop per turn).
MaxConnections=20
# Trigger rate limit made explicit.  An Accept=yes socket that is triggered more
# than TriggerLimitBurst times within TriggerLimitIntervalSec is DISABLED until a
# manual `systemctl restart`.  Tripping it is fail-SAFE (the bridge stops, no
# privileged action runs, the agent's files simply stay SANDBOX_USER-owned and the
# hooks no-op via `|| true`) -- but it would still halt all handback, so the burst
# is set comfortably above any legitimate session's hook rate while still bounding
# a flood.  Raise TriggerLimitBurst if a heavy interactive session ever trips it.
TriggerLimitIntervalSec=10s
TriggerLimitBurst=200
RuntimeDirectory=ai-tools
# Mode 0711, NOT 0750+group.  The runtime dir's GROUP cannot be set declaratively on
# systemd 252 (RuntimeDirectoryGroup= first exists in v256; see the header) and
# RuntimeDirectory= recreates the dir root:root on every (re)start, re-asserting that
# ownership after any ExecStartPost chgrp -- so a chgrp does not stick (confirmed: the
# dir stayed root:root with zero AVCs even under `semodule -DB`, i.e. not an SELinux
# denial).  Mode is the only ownership knob systemd applies reliably here, so the dir
# is made world-TRAVERSABLE (--x) but not listable: any uid can resolve the socket
# name to connect(), nobody but root can enumerate the directory, and the actual
# access gate is unchanged -- the socket file itself is 0660 root:@SANDBOX_GROUP@
# (SocketMode/SocketGroup above) plus the daemon's SO_PEERCRED uid check.  A
# non-@SANDBOX_GROUP@ uid that traverses in still gets EACCES on connect().  This is
# the standard pattern for a group-gated runtime socket on pre-256 systemd.
RuntimeDirectoryMode=0711
ExecStartPost=-/usr/sbin/restorecon -RF /run/ai-tools

[Install]
WantedBy=sockets.target
