# src/usr/lib/systemd/system/ai-tools-relabel.service
[Unit]
Description=Relabel the sandbox claude entrypoint to ai_tools_exec_t
Documentation=man:restorecon(8)
# Triggered by ai-tools-relabel.path; also runnable on demand: systemctl start
# ai-tools-relabel.service.

[Service]
Type=oneshot
# Runs the relabel helper as root in the system instance (which holds relabel privilege),
# triggered by ai-tools-relabel.path after a Node bump. Idempotent and best-effort;
# claude-run's fail-closed preflight is the backstop. See .claude/rules/updater.rule.md.
ExecStart=/usr/local/sbin/ai-tools/ai-tools-relabel-entrypoint
